# profile.d/aliasses.sh
alias runserver="poetry run python manage.py runserver 0.0.0.0:8000"
alias makemigrations="poetry run python manage.py makemigrations"
alias migrate="poetry run python manage.py migrate"
alias shell="poetry run python manage.py shell"